# templates/logstash.init.erb
#!/bin/bash
#
# /etc/rc.d/init.d/logstash
#
# Logstash SysV Init Script
# 
# chkconfig: 2345 80 20
#

# Source function library.
. /etc/init.d/functions

LOGSTASH_DIR="<%= @basedir %>"
LOGSTASH_JAR="<%= @basedir %>/logstash-<%= @pkg_version %>.jar"
LOGSTASH_CONF="<%= @basedir %>/<%= @conf_file %>"
LOGSTASH_USER="<%= @user %>"
LOGSTASH_JAVA_ARGS="-Xms<%= @java_maxheap %> -Xmx<%= @java_maxheap %> <%= @java_opts %>"

start() {
  echo -n "Starting Logstash <%= @server_type %>: "
  daemon --user=$LOGSTASH_USER /usr/bin/java $LOGSTASH_JAVA_ARGS -jar $LOGSTASH_JAR agent -f $LOGSTASH_CONF <% if @server_type == 'indexer' %>-- web<% end %> &> /dev/null &
  echo $?
} 

stop() {
  echo "Shutting down Logstash <%= @server_type %>"
  PID=`/bin/ps -C java -o pid,cmd | grep logstash.*jar | awk '{ print $1; };'`
  EXIT_CODE="0"
  if [ $PID ]; then 
    kill $PID || EXIT_CODE="1"
  fi
  return $EXIT_CODE
}

status() { 
  echo -n "Logstash "
  PID=`/bin/ps -C java -o pid,cmd | grep logstash.*jar | awk '{ print $1; };'`
  if [ $PID ]; then
    echo "is running - PID $PID"
    exit 0
  else
    echo "is not running"
    exit 1
  fi
}

case "$1" in
  start)
    start
  ;;
  stop)
    stop
  ;;
  status)
    status
  ;;
  restart)
    stop
    start
  ;;
  reload)
    echo "reload not yet implemented"
  ;;
    *)
  echo "Usage: <servicename> {start|stop|status|restart}"
  exit 1
  ;;
esac
exit $?
